# systemd/pi5-assistant.service
[Unit]
Description=Pi5 Swedish Voice Assistant (FastAPI)
After=network.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/pi/apps/GENIO
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/apps/GENIO/.venv/bin/uvicorn backend.app:app --host 0.0.0.0 --port 8080
User=pi
Restart=on-failure

[Install]
WantedBy=multi-user.target
